# fetch-crl.cron
#!/bin/sh

if [ -f /etc/default/fetch-crl ] ; then
. /etc/default/fetch-crl
fi

CRLDIR=${CRLDIR:-/etc/grid-security/certificates} 

QUIET=${QUIET:-yes}

if [ $QUIET = "yes" ] ; then
  /usr/sbin/fetch-crl  --loc $CRLDIR --out $CRLDIR -a 24 --quiet
else
  /usr/sbin/fetch-crl  --loc $CRLDIR --out $CRLDIR 
fi

